# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd_err.log
priority=10

[program:bore-ssh]
; FIXED: correct bore syntax is: bore local PORT --to SERVER:PORT
command=bore local 22 --to %(ENV_BORE_SERVER)s:%(ENV_BORE_PORT)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bore-ssh.log
stderr_logfile=/var/log/supervisor/bore-ssh_err.log
priority=20
startretries=10
startsecs=5